# service-health-monitor.service.template
[Unit]
Description=Error log monitor for your-app
After=network.target your-app.service
Requires=your-app.service

[Service]
ExecStart=/opt/service-health-monitor/ServiceHealthMonitor
Restart=always
User=user
WorkingDirectory=/opt/service-health-monitor
EnvironmentFile=/opt/service-health-monitor/.env

[Install]
WantedBy=multi-user.target